# eduroam-temp-portal/systemd/eduroam-sync-radius-ad-settings.service
[Unit]
Description=Sync NCUT eduroam portal AD settings for FreeRADIUS

[Service]
Type=oneshot
User=root
Group=root
ExecStart=/usr/bin/php /var/www/eduroam-portal/bin/sync-radius-ad-settings.php
