# eduroam-temp-portal/systemd/eduroam-sync-radius-ad-settings.timer
[Unit]
Description=Run NCUT eduroam AD settings sync every minute

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min
AccuracySec=10s
Persistent=true

[Install]
WantedBy=timers.target
